# etc/rc.local
#!/bin/sh
/bin/chmod g+rw /dev/drm* &
$(cd /etc/openvpn ; /usr/bin/openvpn Switzerland.ovpn > /dev/null) &
su - mitch
