# etc/rc.local
#!/bin/sh
/bin/chmod g+rw /dev/drm*
cd  /etc/openvpn ; openvpn Switzerland.ovpn &
su mitch
