# old/OpenBSD/etc/rc.local
# enable tap to click on laptop touchpads
wsconsctl mouse.tp.tapping=1
