# openbsd/etc/rc.local
#!/bin/sh

# stop .core file generation on crashes
ulimit -c 0

