# old/crux-etc/rc.local
#!/bin/sh

modprobe fuse tun

ifconfig `/bin/sh /home/mitch/bin/interface` up > /dev/null
wpa_supplicant -i `/bin/sh /home/mitch/bin/interface` -c /etc/wpa_supplicant.conf > /dev/null &

cat > /etc/resolv.conf << EOF
# DNS Server: http://dns.watch
nameserver 84.200.69.80
nameserver 84.200.70.40
EOF

/sbin/sdhcp -d `/bin/sh /home/mitch/bin/interface`

chmod -R 775 /dev/snd*
chmod -R 775 /dev/dri*
chmod -R 777 /dev/shm
